# systemd/system/getty@tty1.service.d/10-getty-tkl-login.conf
[Service]
ExecStartPre=
ExecStartPre=/usr/lib/inithooks/bin/login_script.sh
